# deploy/systemd/nhb.service
[Unit]
Description=NHBChain Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nhb
Group=nhb
WorkingDirectory=/opt/nhbchain
EnvironmentFile=/etc/nhbchain/node.env
ExecStart=/opt/nhbchain/bin/nhb --config /etc/nhbchain/config.toml
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
